# api/config/docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=warn

[program:redis]
command=redis-server --port 6379 --bind 127.0.0.1 --appendonly no --maxmemory 256mb --maxmemory-policy allkeys-lru --save ""
user=rails
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:sidekiq]
command=bundle exec sidekiq -C config/sidekiq.yml
directory=/rails
user=rails
autostart=true
autorestart=true
priority=20
startsecs=10
stopwaitsecs=20
stopsignal=TSTP
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=REDIS_URL="redis://127.0.0.1:6379/0",HOME="/home/rails",PGSSLCERT="",PGSSLKEY=""

[program:puma]
command=bundle exec rails server -b 0.0.0.0 -p 3000
directory=/rails
user=rails
autostart=true
autorestart=true
priority=30
startsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=REDIS_URL="redis://127.0.0.1:6379/0",HOME="/home/rails",PGSSLCERT="",PGSSLKEY=""
